# roles/docker-api/templates/thuv-docker.service.j2
[Unit]
Description="Thuv docker API"

[Service]
ExecStart={{ web_docker_api_location }}docker-api

[Install]
WantedBy=multi-user.target
